# config/crontab
* * * * * (cd ~/github/dkorolev/homepage; git pull)
